# k3s/rootfs/etc/s6-overlay/s6-rc.d/k3s/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: NK3S8N
# Runs K3S script
# ==============================================================================

bashio::log.info "Starting K3S..."

exec k3s server --debug --config /data/config.yaml