# k3s/rootfs/etc/s6-overlay/s6-rc.d/k3s/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: K3S
# Runs K3S script
# ==============================================================================

bashio::log.info "Starting K3S..."

# Ensure configuration data directory exists
if ! bashio::fs.directory_exists '/config/data'; then
    mkdir -p /config/data || bashio::exit.nok "Failed to create k3s configuration directory"
fi

exec k3s server --debug --write-kubeconfig-mode "0644" --data-dir "/config/data"