# k3s/rootfs/etc/s6-overlay/s6-rc.d/init-k3s/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: K3S
# K3S init script, runs before any other service
# ==============================================================================

bashio::log.info "Configuring K3S..."

declare db_url

# Ensure configuration data directory exists
if ! bashio::fs.directory_exists '/config/data'; then
    mkdir -p /config/data || bashio::exit.nok "Failed to create k3s configuration directory"
fi

if ["$(bashio::config.true 'db_type')" == "mysqldb"] == 'mysqldb'; then
  db_url="mysql://$(bashio::config 'db_user'):$(bashio::config 'db_password')@tcp($(bashio::config 'db_host'):^$(bashio::config 'db_port'))/$(bashio::config 'db_name')";
elif ["$(bashio::config.true 'db_type')" == "postgresdb"]; then
  db_url="postgres://$(bashio::config 'db_user'):$(bashio::config 'db_password')@$(bashio::config 'db_host'):^$(bashio::config 'db_port')/$(bashio::config 'db_name')";
else
  db_url="file:/config/data/database.db"
fi

bashio::var.json \
    db_url "$db_url" \
    | tempio \
        -template /etc/k3s/config.gtpl \
        -out /data/config.yaml

if [ "$(bashio::config 'log_level')" == "debug" ]; then
  bashio::log.info "$(cat /data/config.yaml)"
fi

